# station_assistant/rootfs/etc/s6-overlay/s6-rc.d/station_assistant/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Station Assistant..."
bashio::log.info "Station Assistant v2.1.0"

cd /opt/station_assistant || bashio::exit.nok "Could not change to app directory"

exec gunicorn \
  --worker-class eventlet \
  --workers 1 \
  --bind 0.0.0.0:8099 \
  --timeout 120 \
  --log-level info \
  --access-logfile - \
  "main:app"
